# systemd/atomic-compose-rawhide.service
[Unit]

[Service]
WorkingDirectory=/srv/fedora-atomic/rawhide
ExecStart=/usr/bin/rpm-ostree-toolbox taskrunner /srv/inbox/rawhide

[Install]
WantedBy=multi-user.target
