# systemd/fedmsg-atomic-composer.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/fedmsg-hub
Type=simple
User=composer
Group=composer

[Install]
WantedBy=multi-user.target
